# Scripts/com.myasu.script.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>ANY_LABEL_STRING-com.myasu.script.plist</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python</string>
        <string>PATH_TO_SCRIPT</string>
    </array>
    <!-- Determines the script running interval -->
    <key>StartInterval</key>
        <integer>300</integer>
    <key>StandardOutPath</key>
	<string>PATH_TO_OUTPUT_FILE</string>
    <key>StandardErrorPath</key>
	<string>PATH_TO_ERROR_OUT_FILE</string>
    <key>RunAtLoad</key>
    <true/>
</dict>
</plist>
